# samples/dserver-update-keycache.timer.sample
[Unit]
Description=Refresh DServer SSH Key Cache every 30 minutes.

[Timer]
OnCalendar=*:0/30
