# samples/dserver-update-keycache.service.sample
[Unit]
Description=Refresh DServer SSH Key Cache

[Service]
Type=oneshot
ExecStart=/var/run/dserver/update_key_cache.sh
